# deploy/keldurben-server.service
[Unit]
Description=KeldurbenHub Game Server
After=network.target

[Service]
Type=simple
User=keldurben
Group=keldurben
WorkingDirectory=/opt/keldurben-server
ExecStart=/opt/keldurben-server/server
Restart=always
RestartSec=5
Environment=RUST_LOG=info
Environment=JWT_SECRET=your-production-secret-key-change-this
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/keldurben-server

[Install]
WantedBy=multi-user.target
